# backend/cron/mail_to_ticket.cron
5 * * * * /path/to/your/opencop/installation/backend/mail_to_ticket.pl